# ops/syrmos-api/systemd/syrmos-backup.timer
[Unit]
Description=Daily SQLite backup for Syrmos

[Timer]
OnCalendar=*-*-* 03:30:00 Europe/Athens
Persistent=true
Unit=syrmos-backup.service

[Install]
WantedBy=timers.target
