# ops/syrmos-api/systemd/syrmos-watcher-oasa.timer
[Unit]
Description=Daily Syrmos OASA fares HTML watcher
Requires=syrmos-watcher-oasa.service

[Timer]
# Staggered against the other watchers (Hellenic Train 03:17, STASY 03:42)
# so DNS/upstream stalls only ever block one feed at a time.
OnCalendar=*-*-* 04:07:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
