# ops/syrmos-api/systemd/syrmos-scraper-24mmm.service
[Unit]
Description=OASA 24mmm scraper one-shot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=peterdsp
WorkingDirectory=/home/peterdsp/syrmos-api
Environment=SYRMOS_DB_PATH=/home/peterdsp/syrmos-api/db/syrmos.db
Environment=SYRMOS_API_OUT_DIR=/home/peterdsp/syrmos-api/out
ExecStart=/home/peterdsp/syrmos-api/.venv/bin/python -m syrmos_admin.scraper_24mmm
ExecStartPost=/home/peterdsp/syrmos-api/.venv/bin/python -m syrmos_admin.generator
StandardOutput=append:/home/peterdsp/syrmos-api/scraper.log
StandardError=append:/home/peterdsp/syrmos-api/scraper.log
